# clearstack.service
[Unit]
Description=Clearstack Daemon
Documentation=https://clearstackhq.com/docs
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStartPre=/bin/mkdir -p /tmp/clearstack
ExecStartPre=/bin/mkdir -p /var/log
ExecStart=/usr/local/bin/clearstack daemon --host 0.0.0.0 --port 4242 --max-workers 10
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clearstack

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp/clearstack /var/log
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target 