# clearly.service
[Unit]
Description=Clearly Daemon
Documentation=https://clearly.run/docs
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStartPre=/bin/mkdir -p /tmp/clearly
ExecStartPre=/bin/mkdir -p /var/log
ExecStart=/usr/local/bin/clearly daemon --host 0.0.0.0 --port 4242 --max-workers 10
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clearly

# Security settings
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths=/tmp/clearly /var/log
# ProtectKernelTunables=true
# ProtectKernelModules=true
# ProtectControlGroups=true
# RestrictRealtime=true
# RestrictSUIDSGID=true
# NoNewPrivileges=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target 